# docker/etc/cron.d/mariadb-backup-aws
SHELL=/bin/sh
0 1 * * *   /app/bin/console mariaDbBackup:aws >> /app/log/cron.log 2>&1
