# docker/etc/supervisor/supervisor.conf
[supervisord]
nodaemon=true

[program:cron]
command=supercronic /etc/cron.d/fingather
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true
